# usr/lib/systemd/user/drella-update-orchestrator.service
[Unit]
Description=Pull latest orchestrator code from GitHub

[Service]
Type=oneshot
ExecStart=/usr/libexec/drella-update-orchestrator
